# test.service
[Unit]
Description=SSU_Crontab daemon service

[Service]
Type=simple
ExecStart=/
Restart=on-failure

[Install]
WantedBy=multi-user.target
